# packages/net-fs/samba/files/systemd/nmbd.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=nmbd server
Before=samba.target

[Service]
Environment=TMP=""
Environment=TMPDIR=""
ExecStart=/usr/sbin/nmbd -F
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
