# packages/net-fs/samba/files/systemd/smbd.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=smbd server
Before=samba.target

[Service]
Environment=TMP=""
ExecStart=/usr/sbin/smbd -F
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
