# packages/sys-apps/pcsc-lite/files/systemd/pcscd.service
[Unit]
Description=pcscd smart card reader daemon

[Service]
User=pcscd
Group=pcscd
ExecStart=/usr/sbin/pcscd --foreground --auto-exit
ExecReload=/usr/sbin/pcscd --hotplug
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
