# packages/sys-apps/pcsc-lite/files/systemd/pcscd.socket
[Unit]
Description=PC/SC Smart Card Daemon Activation Socket

[Socket]
User=pcscd
Group=pcscd
ListenStream=/run/pcscd/pcscd.comm

[Install]
WantedBy=sockets.target
